# debian/python3-cyberfusion-external-providers-ip-ranges.cron.d
25 18 * * * root /usr/bin/tuxis-cron /etc/cronscripts/update-ip-ranges-ferm-variables.cron
